# systemd/rclone-mount@.service
[Unit]
Description=Rclone mount for %i cloud storage
After=network-online.target
Wants=network-online.target
Before=multi-stream.service

[Service]
Type=simple
User=root
Group=root
Environment=RCLONE_CONFIG=/root/.config/rclone/rclone.conf
ExecStartPre=/bin/mkdir -p /mnt/%i
ExecStart=/usr/bin/rclone mount %i: /mnt/%i \
  --config=/root/.config/rclone/rclone.conf \
  --allow-other \
  --vfs-cache-mode writes \
  --vfs-cache-max-size 2G \
  --vfs-cache-max-age 1h \
  --vfs-read-chunk-size 64M \
  --vfs-read-chunk-size-limit 2G \
  --buffer-size 64M \
  --dir-cache-time 1h \
  --umask 002 \
  --log-level INFO \
  --log-file /var/log/rclone-%i.log
ExecStop=/bin/fusermount -uz /mnt/%i
ExecStopPost=/bin/rmdir /mnt/%i
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
KillMode=process

[Install]
WantedBy=multi-user.target
