# systemd/stream-cleanup.timer
[Unit]
Description=Run stream cleanup service periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target